# deb_build/dualrecord_1.0.0_amd64/DEBIAN/postinst
#!/bin/bash
# 安装后脚本：注册 systemd 服务、配置 XDG 自启动、立即启动

APP_NAME="dualrecord"
SERVICE_FILE="/etc/systemd/system/${APP_NAME}.service"
INSTALL_BIN="/opt/dualrecord/bin/DualRecordService"

# 获取实际桌面登录用户
get_desktop_user() {
    local u=""
    u=$(logname 2>/dev/null) || true
    [ -z "$u" ] || [ "$u" = "root" ] && u="${SUDO_USER}"
    [ -z "$u" ] || [ "$u" = "root" ] && u=$(who 2>/dev/null | grep -v root | awk 'NR==1{print $1}') || true
    [ -z "$u" ] && u="${USER}"
    echo "$u"
}

CURRENT_USER=$(get_desktop_user)
if [ -z "${CURRENT_USER}" ] || [ "${CURRENT_USER}" = "root" ]; then
    echo "[WARN] 无法确定桌面用户，服务将以 root 运行（不推荐）"
    CURRENT_USER="root"
fi

USER_HOME=$(getent passwd "${CURRENT_USER}" 2>/dev/null | cut -d: -f6)
[ -z "${USER_HOME}" ] && USER_HOME=$(eval echo ~${CURRENT_USER} 2>/dev/null)
[ -z "${USER_HOME}" ] && USER_HOME="/home/${CURRENT_USER}"

echo "[INFO] 配置用户: ${CURRENT_USER}  主目录: ${USER_HOME}"

# 替换 service 文件中的占位符
sed -i "s|DUALRECORD_USER|${CURRENT_USER}|g"   "${SERVICE_FILE}"
sed -i "s|DUALRECORD_HOME|${USER_HOME}|g"       "${SERVICE_FILE}"

# 创建用户工作目录
if [ "${CURRENT_USER}" != "root" ]; then
    mkdir -p "${USER_HOME}/BankDualRecord/videos"
    mkdir -p "${USER_HOME}/BankDualRecord/logs"
    chown -R "${CURRENT_USER}:${CURRENT_USER}" "${USER_HOME}/BankDualRecord" 2>/dev/null || true
fi

# XDG 自启动（兜底，在没有 systemd 图形会话的情况下也能自动启动）
if [ "${CURRENT_USER}" != "root" ] && [ -d "${USER_HOME}" ]; then
    AUTOSTART_DIR="${USER_HOME}/.config/autostart"
    mkdir -p "${AUTOSTART_DIR}"
    cat > "${AUTOSTART_DIR}/dualrecord.desktop" << EOF
[Desktop Entry]
Type=Application
Name=银行双录控件
Exec=${INSTALL_BIN}
Icon=dualrecord
StartupWMClass=BankDualRecord
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
EOF
    chown -R "${CURRENT_USER}:${CURRENT_USER}" "${AUTOSTART_DIR}" 2>/dev/null || true
fi

# 注册并启动服务
systemctl daemon-reload 2>/dev/null || true
systemctl enable "${APP_NAME}.service" 2>/dev/null || true

# 停止旧进程（支持重新安装场景）
systemctl stop "${APP_NAME}.service" 2>/dev/null || true
pkill -9 -f "DualRecordService" 2>/dev/null || true
sleep 1

# 尝试 systemctl start
systemctl start "${APP_NAME}.service" 2>/dev/null
START_RC=$?

if [ ${START_RC} -eq 0 ]; then
    sleep 1
    if systemctl is-active --quiet "${APP_NAME}.service" 2>/dev/null; then
        echo "[INFO] 服务启动成功 ✓"
    else
        echo "[WARN] 服务已提交启动，正在初始化..."
        START_RC=1
    fi
fi

# fallback：直接以用户身份拉起进程
if [ ${START_RC} -ne 0 ]; then
    echo "[WARN] systemctl 启动未就绪，尝试直接启动进程..."
    DISPLAY_VAR=$(su - "${CURRENT_USER}" -c 'echo $DISPLAY' 2>/dev/null || echo ":0")
    nohup su - "${CURRENT_USER}" -c \
        "DISPLAY=${DISPLAY_VAR} XAUTHORITY=${USER_HOME}/.Xauthority ${INSTALL_BIN}" \
        >/tmp/dualrecord_start.log 2>&1 &
    sleep 2
    if pgrep -f "DualRecordService" >/dev/null 2>&1; then
        echo "[INFO] 服务直接启动成功 ✓"
    else
        echo "[WARN] 即时启动未成功，开机登录后将自动运行"
    fi
fi

# 更新图标缓存
gtk-update-icon-cache /usr/share/icons/hicolor 2>/dev/null || true
update-desktop-database 2>/dev/null || true

echo ""
echo "================================================"
echo "  银行柜面双录控件 安装完成！"
echo "================================================"
echo "  程序路径: /opt/dualrecord/bin/DualRecordService"
echo "  HTTP端口: 9528 (可在初始化向导中修改)"
echo "  视频目录: ${USER_HOME}/BankDualRecord/videos"
echo "  查看日志: journalctl -u dualrecord -f"
echo "================================================"
