# deb_build/dualrecord_1.0.0_amd64/DEBIAN/prerm
#!/bin/bash
# 卸载前脚本：停止服务
APP_NAME="dualrecord"
echo "[INFO] 停止双录服务..."
systemctl stop "${APP_NAME}.service"   2>/dev/null || true
systemctl disable "${APP_NAME}.service" 2>/dev/null || true
pkill -9 -f "DualRecordService"         2>/dev/null || true
